# deploy/systemd/tradebot-dashboard.service
[Unit]
Description=TradeBot mission control dashboard (read-only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tradebot
Group=tradebot
WorkingDirectory=/opt/tradebot
EnvironmentFile=/opt/tradebot/.env

# Binds to 0.0.0.0 so Tailscale-connected devices (PC, phone) can reach it
# directly at http://100.113.140.69:8080. UFW has no allow rule for 8080 so
# the public internet cannot reach it. Never add a UFW allow for 8080 without
# adding token auth first (Phase 3).
Environment=DASHBOARD_HOST=0.0.0.0
Environment=DASHBOARD_PORT=8080
# DASHBOARD_TOKEN gates POST /api/bot/restart and POST /api/bot/stop.
# Set it in /opt/tradebot/.env (loaded via EnvironmentFile above) — never
# hardcode it here. If unset, control endpoints return 503 (fail closed).
# Also requires /etc/sudoers.d/tradebot-dashboard (see deploy/sudoers/).

ExecStart=/opt/tradebot/venv/bin/python -m dashboard

Restart=on-failure
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tradebot-dashboard

[Install]
WantedBy=multi-user.target
