# deploy/systemd/tradebot-notify@.service
[Unit]
Description=Failure notification for %i
# Template unit — instantiated as tradebot-notify@<service-name>.service

[Service]
Type=oneshot
# Collect last 50 journal lines from the failed service, then push to ntfy.sh
ExecStart=/bin/bash -c '\
    SERVICE="%i"; \
    LOGS=$(journalctl -u "$SERVICE" -n 50 --no-pager --output=short 2>/dev/null); \
    BODY="$SERVICE failed at $(date -u +%%Y-%%m-%%d\ %%H:%%M\ UTC)\n\n--- last 50 log lines ---\n${LOGS}"; \
    curl -s --max-time 10 \
        -X POST "https://ntfy.sh/tradebot-DUE090987" \
        -H "Title: TradeBot FAILURE: $SERVICE" \
        -H "Priority: urgent" \
        -H "Tags: warning,rotating_light" \
        -d "$BODY" \
    || true'
