# deploy/systemd/tradebot-notify@.service
[Unit]
Description=Failure notification for %i
# Template unit — instantiated as tradebot-notify@<service-name>.service

[Service]
Type=oneshot
EnvironmentFile=/opt/tradebot/.env
# Send a content-free failure alert — logs stay on the VPS (check with journalctl)
ExecStart=/bin/bash -c '\
    SERVICE="%i"; \
    BODY="$SERVICE failed at $(date -u +%%Y-%%m-%%d\ %%H:%%M\ UTC). Check VPS: journalctl -u $SERVICE -n 50 --no-pager"; \
    curl -s --max-time 10 \
        -X POST "https://ntfy.sh/${NTFY_TOPIC}" \
        -H "Title: TradeBot FAILURE: $SERVICE" \
        -H "Priority: urgent" \
        -H "Tags: warning,rotating_light" \
        -d "$BODY" \
    || true'
