# deploy/systemd/tradebot-dashboard.service
[Unit]
Description=TradeBot mission control dashboard (read-only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tradebot
Group=tradebot
WorkingDirectory=/opt/tradebot
EnvironmentFile=/opt/tradebot/.env

# Bind directly to the Tailscale interface — PC + iPhone reach it via
# http://100.113.140.69:8080, and there is no public listener at all.
# UFW already blocks 8080 on the public NIC, but binding only to the
# Tailscale IP is defense-in-depth: even if UFW is misconfigured the
# socket cannot accept on 0.0.0.0.
# If the Tailscale IP ever changes, update this line and `daemon-reload`.
Environment=DASHBOARD_HOST=100.113.140.69
Environment=DASHBOARD_PORT=8080
# DASHBOARD_TOKEN gates POST /api/bot/restart and POST /api/bot/stop.
# Set it in /opt/tradebot/.env (loaded via EnvironmentFile above) — never
# hardcode it here. If unset, control endpoints return 503 (fail closed).
# Also requires /etc/sudoers.d/tradebot-dashboard (see deploy/sudoers/).

ExecStart=/opt/tradebot/venv/bin/python -m dashboard

Restart=on-failure
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=tradebot-dashboard

[Install]
WantedBy=multi-user.target
