# deploy/systemd/tradebot-health.service
[Unit]
Description=TradeBot daily health check

[Service]
Type=oneshot
ExecStart=/bin/bash -c '\
    HEALTH_FILE=/opt/tradebot/data/health.txt; \
    NTFY_TOPIC=tradebot-DUE090987; \
    if [ ! -f "$HEALTH_FILE" ]; then \
        curl -s --max-time 10 \
            -X POST "https://ntfy.sh/$NTFY_TOPIC" \
            -H "Title: TradeBot health: health.txt missing" \
            -H "Priority: high" \
            -H "Tags: warning" \
            -d "Health file not found at $(date -u). Bot may not have run since deploy." \
        || true; \
        exit 0; \
    fi; \
    LAST=$(cat "$HEALTH_FILE"); \
    NOW=$(date -u +%%s); \
    LAST_TS=$(date -u -d "$LAST" +%%s 2>/dev/null || echo 0); \
    AGE=$(( NOW - LAST_TS )); \
    if [ "$AGE" -gt 93600 ]; then \
        curl -s --max-time 10 \
            -X POST "https://ntfy.sh/$NTFY_TOPIC" \
            -H "Title: TradeBot health: on_tick stale" \
            -H "Priority: high" \
            -H "Tags: warning" \
            -d "Last tick was at $LAST ($(( AGE / 3600 ))h ago). Bot may be stuck or not firing at 16:10 ET." \
        || true; \
    fi'
